# com.claude-usage-menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-usage-menubar</string>
    <key>ProgramArguments</key>
    <array>
        <!--
            Update these two paths to match your system:
            1. Your python3 path (run: which python3)
            2. Where you cloned this repo
        -->
        <string>/opt/homebrew/bin/python3</string>
        <string>/path/to/claude-usage-menubar/claude_menubar.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude_menubar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude_menubar.err</string>
</dict>
</plist>
